# backups/backup-casey.cron
# dot-tags hostname:casey

10 */2 * * * bash -l -c "~/backup-scripts/backup-casey-to-borg.sh >> ~/backup-scripts/.log 2>&1"
