# backups/backup-casey.cron
# dot-tags hostname:casey

10 */2 * * * bash -l -c "cd ~/backup-scripts; source .secrets; ./backup-casey-to-morgan.sh >> .log 2>&1"
