# backups/backup-morgan.cron
# dot-tags hostname:morgan

30 */2 * * * bash -l -c "cd ~/backup-scripts; source .secrets; ./backup-archive-to-morgan.sh >> .log 2>&1"
40 3 * * * bash -l -c "cd ~/backup-scripts; source .secrets; ./prune-morgan-backups.sh >> .log 2>&1"
50 4 * * * bash -l -c "cd ~/backup-scripts; source .secrets; ./sync-morgan-to-b2.sh >> .log 2>&1"
00 9 * * * bash -l -c "cd ~/backup-scripts; source .secrets; initnode; node ./send-backup-report.js >> .log 2>&1"
