# photos/sync-wallpapers.cron
# dot-tags hostname:casey

30 * * * * bash -l -c "node ~/Pictures/.utils/sync-wallpapers.js"
