# vim/vimwiki/convert-to-html.cron
# dot-tags hostname:chuck

2-59/5 * * * * bash -l -c "cd ~/dotfiles; ./vim/vimwiki/convert-to-html.sh"
