# vim/vimwiki/create-diary-entry.cron
# dot-tags hostname:chuck

10 0 * * * bash -l -c "cd ~/dotfiles; ./vim/vimwiki/create-diary-entry.sh"
