# vim/vimwiki/update-tag-pages.cron
# dot-tags hostname:chuck

2/5 * * * * bash -l -c "cd ~/dotfiles; ./vim/vimwiki/update-tag-pages.sh"
