# dotfiles-meta/dot-sync.cron
*/15 * * * * bash -l -c "cd ~/dotfiles; ./bin/dot-sync &>> ./dotfiles-meta/.cron-log"
