# dotfiles-meta/dot-sync.cron
# dot-tags profile:personal

* * * * * bash -l -c "cd ~/dotfiles; ./dotfiles-meta/self-update.sh &>> ./dotfiles-meta/.cron-log; ./bin/dot-sync &>> ./dotfiles-meta/.cron-log"
